# startup/mosesMT.enfr.service
[Unit]
Description=Moses MT Web Service (EN-FR)
After=multi-user.target

[Service]
ExecStart=/home/ubuntu/mosesdecoder/bin/moses --server --server-port 3004 -f /home/ubuntu/mt_publish/model_enfr_ted/moses.ini --mark-unknown --unknown-word-prefix "__UNK__" --unknown-word-suffix "__UNK__"
Restart=on-failure

[Install]
WantedBy=default.target


