# startup/mosesMT.enzh.service
[Unit]
Description=Moses MT Web Service (EN-ZH)
After=multi-user.target

[Service]
ExecStart=/home/ubuntu/mosesdecoder/bin/moses --server --server-port 3001 -f /home/ubuntu/mt_publish/model_enzh_ted/moses.ini --mark-unknown --unknown-word-prefix "__UNK__" --unknown-word-suffix "__UNK__"
Restart=on-failure

[Install]
WantedBy=default.target


